# ubuntu-rootfs/etc/rc.local
#!/bin/bash

#config your app autorun there when sytem boot

#4G-demon &
