# ubuntu-rootfs/usr/lib/systemd/user/indicators-pre.target
[Unit]
Description=Target representing the lifecycle of the indicators. Each indicator should be bound to it in its individual service file.
