# ubuntu-rootfs/var/lib/dpkg/info/bind9.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/force-complain
/etc/apparmor.d/local
/etc/apparmor.d/usr.sbin.named
/etc/bind
/etc/bind/bind.keys
/etc/bind/db.0
/etc/bind/db.127
/etc/bind/db.255
/etc/bind/db.empty
/etc/bind/db.local
/etc/bind/db.root
/etc/bind/named.conf
/etc/bind/named.conf.default-zones
/etc/bind/named.conf.local
/etc/bind/named.conf.options
/etc/bind/zones.rfc1918
/etc/init.d
/etc/init.d/bind9
/etc/insserv.conf.d
/etc/insserv.conf.d/bind9
/etc/network
/etc/network/if-down.d
/etc/network/if-down.d/bind9
/etc/network/if-up.d
/etc/network/if-up.d/bind9
/etc/ppp
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/bind9
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/bind9
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/bind9
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bind9-pkcs11.service
/lib/systemd/system/bind9-resolvconf.service
/lib/systemd/system/bind9.service
/usr
/usr/bin
/usr/bin/arpaname
/usr/bin/bind9-config
/usr/bin/named-rrchecker
/usr/lib
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/bind9.conf
/usr/sbin
/usr/sbin/ddns-confgen
/usr/sbin/dnssec-importkey
/usr/sbin/genrandom
/usr/sbin/isc-hmac-fixup
/usr/sbin/named
/usr/sbin/named-journalprint
/usr/sbin/named-pkcs11
/usr/sbin/nsec3hash
/usr/share
/usr/share/bind9
/usr/share/doc
/usr/share/doc/bind9
/usr/share/doc/bind9/NEWS.Debian.gz
/usr/share/doc/bind9/README.Debian.gz
/usr/share/doc/bind9/README.gz
/usr/share/doc/bind9/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/arpaname.1.gz
/usr/share/man/man1/bind9-config.1.gz
/usr/share/man/man1/isc-config.sh.1.gz
/usr/share/man/man1/named-rrchecker.1.gz
/usr/share/man/man5
/usr/share/man/man5/named.conf.5.gz
/usr/share/man/man5/rndc.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/ddns-confgen.8.gz
/usr/share/man/man8/dnssec-coverage.8.gz
/usr/share/man/man8/dnssec-importkey.8.gz
/usr/share/man/man8/genrandom.8.gz
/usr/share/man/man8/isc-hmac-fixup.8.gz
/usr/share/man/man8/named-journalprint.8.gz
/usr/share/man/man8/named.8.gz
/usr/share/man/man8/nsec3hash.8.gz
/var
/var/cache
/var/cache/bind
/usr/sbin/tsig-keygen
/usr/share/doc/bind9/changelog.Debian.gz
/usr/share/man/man8/tsig-keygen.8.gz
