# ubuntu-rootfs/var/lib/dpkg/info/bind9.conffiles
/etc/apparmor.d/usr.sbin.named
/etc/bind/bind.keys
/etc/bind/db.0
/etc/bind/db.127
/etc/bind/db.255
/etc/bind/db.empty
/etc/bind/db.local
/etc/bind/db.root
/etc/bind/named.conf
/etc/bind/named.conf.default-zones
/etc/bind/named.conf.local
/etc/bind/named.conf.options
/etc/bind/zones.rfc1918
/etc/init.d/bind9
/etc/insserv.conf.d/bind9
/etc/network/if-down.d/bind9
/etc/network/if-up.d/bind9
/etc/ppp/ip-down.d/bind9
/etc/ppp/ip-up.d/bind9
/etc/ufw/applications.d/bind9
