# ubuntu-rootfs/var/lib/dpkg/info/containerd.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/containerd.service
/usr
/usr/bin
/usr/bin/containerd
/usr/bin/containerd-shim
/usr/bin/containerd-shim-runc-v1
/usr/bin/containerd-shim-runc-v2
/usr/bin/containerd-stress
/usr/bin/ctr
/usr/share
/usr/share/doc
/usr/share/doc/containerd
/usr/share/doc/containerd/PLUGINS.md.gz
/usr/share/doc/containerd/README.md.gz
/usr/share/doc/containerd/RUNC.md
/usr/share/doc/containerd/SECURITY_AUDIT.pdf.gz
/usr/share/doc/containerd/changelog.Debian.gz
/usr/share/doc/containerd/client-opts.md
/usr/share/doc/containerd/content-flow.md.gz
/usr/share/doc/containerd/copyright
/usr/share/doc/containerd/garbage-collection.md.gz
/usr/share/doc/containerd/getting-started.md.gz
/usr/share/doc/containerd/hosts.md.gz
/usr/share/doc/containerd/managed-opt.md.gz
/usr/share/doc/containerd/namespaces.md
/usr/share/doc/containerd/ops.md.gz
/usr/share/doc/containerd/remote-snapshotter.md.gz
/usr/share/doc/containerd/rootless.md
/usr/share/doc/containerd/stream_processors.md
/usr/share/doc/containerd/tracing.md
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/containerd-config.toml.5.gz
/usr/share/man/man8
/usr/share/man/man8/containerd-config.8.gz
/usr/share/man/man8/containerd.8.gz
/usr/share/man/man8/ctr.8.gz
