# ubuntu-rootfs/var/lib/dpkg/info/isc-dhcp-server.conffiles
/etc/apparmor.d/usr.sbin.dhcpd
/etc/dhcp/dhcpd.conf
/etc/dhcp/dhcpd6.conf
/etc/init.d/isc-dhcp-server
/etc/logcheck/ignore.d.server/isc-dhcp-server
