# ubuntu-rootfs/var/lib/dpkg/info/netfilter-persistent.list
/.
/etc
/etc/default
/etc/default/netfilter-persistent
/etc/init.d
/etc/init.d/netfilter-persistent
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/netfilter-persistent.service
/usr
/usr/sbin
/usr/sbin/netfilter-persistent
/usr/share
/usr/share/doc
/usr/share/doc/netfilter-persistent
/usr/share/doc/netfilter-persistent/README
/usr/share/doc/netfilter-persistent/changelog.gz
/usr/share/doc/netfilter-persistent/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/netfilter-persistent.8.gz
/usr/share/netfilter-persistent
/usr/share/netfilter-persistent/plugins.d
