# ubuntu-rootfs/var/lib/dpkg/info/netfilter-persistent.conffiles
/etc/default/netfilter-persistent
/etc/init.d/netfilter-persistent
