# ubuntu-rootfs/var/lib/dpkg/info/ubuntu-fan.list
/.
/etc
/etc/NetworkManager
/etc/NetworkManager/conf.d
/etc/NetworkManager/conf.d/10-ubuntu-fan.conf
/etc/default
/etc/default/ubuntu-fan
/etc/dnsmasq.d
/etc/dnsmasq.d/ubuntu-fan
/etc/init
/etc/init/ubuntu-fan.conf
/etc/init.d
/etc/init.d/ubuntu-fan
/etc/network
/etc/network/if-down.d
/etc/network/if-down.d/ubuntu-fan
/etc/network/if-post-down.d
/etc/network/if-post-down.d/ubuntu-fan
/etc/network/if-up.d
/etc/network/if-up.d/ubuntu-fan
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ubuntu-fan.service
/usr
/usr/lib
/usr/lib/ubuntu-fan
/usr/lib/ubuntu-fan/fan-net
/usr/sbin
/usr/sbin/fanatic
/usr/sbin/fanctl
/usr/share
/usr/share/doc
/usr/share/doc/ubuntu-fan
/usr/share/doc/ubuntu-fan/README.gz
/usr/share/doc/ubuntu-fan/changelog.gz
/usr/share/doc/ubuntu-fan/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/fanatic.8.gz
/usr/share/man/man8/fanctl.8.gz
/var
/var/lib
/var/lib/ubuntu-fan
/var/lib/ubuntu-fan/config
/var/lib/ubuntu-fan/config/README
